# Ansible/startup.service
[Unit]
After=network.service

[Service]
ExecStart=/usr/local/bin/run.sh

[Install]
WantedBy=default.target
