# Ansible/startupslave.service
[Unit]
After=network.service

[Service]
ExecStart=/usr/local/bin/runslave.sh
Restart=on-failure

[Install]
WantedBy=default.target


